# 06-systemd/timer/auto_clean.timer
[Unit]
Description=A timer for Auto clean


[Timer]
OnCalendar=monthly
Unit=auto_clean.service
Persistent=true

[Install]
WantedBy=timers.target